# docker/supervisor/task_room_cache_recalculate.conf
[program:task_room_cache_recalculate]
command=/var/www/mbh/bin/console rabbitmq:consumer task_room_cache_recalculate
directory=/var/www/mbh
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor.err.log
stdout_logfile=/var/log/supervisor.out.log
user=www-data